# templates/mastodon-schedule.timer.j2
[Unit]
Description=Mastodon Toot Schedule Timer

[Timer]
OnCalendar={{ mastodon_scheduler__time }}
RandomizedDelaySec=300s
Persistent=yes

[Install]
WantedBy=network.target
